# .config/systemd/user/backup-keepass-git.timer
[Unit]
Description=Run backup KeePassXC script once per day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
