# .config/systemd/user/xset-rate.service
[Unit]
Description=Set X11 Key Repeat Rate
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/xset r rate 275 23
