# .config/systemd/user/xset-rate.timer
[Unit]
Description=Run xset rate hourly

[Timer]
OnCalendar=hourly
Persistent=true
Unit=xset-rate.service

[Install]
WantedBy=timers.target
