# .config/systemd/user/xset-rate.timer
[Unit]
Description=Run xset rate hourly

[Timer]
OnBootSec=0
OnUnitActiveSec=1h
Persistent=true
Unit=xset-rate.service

[Install]
WantedBy=timers.target
